# deploy/systemd/mengasr-listener.service
[Unit]
Description=MengASR2 Listener — HTTP API 网关（对公网暴露）
After=network.target mengasr-worker.service
Requires=mengasr-worker.service

[Service]
Type=simple
User=bchnm
WorkingDirectory=/srv/mengasr/server

# ── 配置文件 ──────────────────────────────────────────────
Environment=MENGASR_CONFIG=/srv/mengasr/mengasr.yaml
Environment=PIP_CONSTRAINT=/srv/mengasr/constraints.txt

ExecStart=/srv/mengasr/.venv/bin/python -m mengasr_server --listener
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mengasr-listener

# 资源限制
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
