# deploy/systemd/mengasr-listener.service.example
[Unit]
Description=MengASR2 Listener — HTTP API 网关
After=network.target
# 根据使用的后端，选择对应的 Worker 服务：
# Wants=mengasr-worker.service        # MiMo 后端
# Wants=mengasr-worker-qwen3.service  # Qwen3 后端

[Service]
Type=simple
# ⚠️ 修改为实际运行用户
User=mengasr
# ⚠️ 修改为实际安装目录
WorkingDirectory=/srv/mengasr/server

# ── 配置文件 ──────────────────────────────────────────────
# ⚠️ 修改为实际配置文件路径
Environment=MENGASR_CONFIG=/srv/mengasr/mengasr.yaml
# ⚠️ 修改为实际约束文件路径
Environment=PIP_CONSTRAINT=/srv/mengasr/constraints.txt

# ⚠️ 修改为实际 Python 解释器路径
ExecStart=/srv/mengasr/.venv/bin/python -m mengasr_server --listener
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mengasr-listener

# 资源限制
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
