# deploy/systemd/mengasr-worker-qwen3.service.example
[Unit]
Description=MengASR2 Worker (Qwen3-ASR) — ASR 模型推理进程
After=network.target
Wants=network.target

[Service]
Type=simple
# ⚠️ 修改为实际运行用户
User=mengasr
# ⚠️ 修改为实际安装目录
WorkingDirectory=/srv/mengasr/server

# ── Qwen3-ASR 独立 venv（transformers 版本与 MiMo 不兼容）───
# ⚠️ 修改 PATH 中的 venv 路径和 CUDA 路径
Environment=PATH=/srv/mengasr/.venv-qwen3asr/bin:/usr/local/cuda/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# ⚠️ 修改为实际配置文件路径
Environment=MENGASR_CONFIG=/srv/mengasr/mengasr-qwen3.yaml
# HuggingFace 镜像（中国大陆无法直连时取消注释）
# Environment=HF_ENDPOINT=https://hf-mirror.com
# ⚠️ 修改为实际约束文件路径
Environment=PIP_CONSTRAINT=/srv/mengasr/constraints-qwen3.txt

# ⚠️ 修改为实际 Python 解释器路径
ExecStart=/srv/mengasr/.venv-qwen3asr/bin/python -m mengasr_server --worker
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mengasr-worker-qwen3

# GPU 访问
SupplementaryGroups=video

# 资源限制
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
