# deploy/systemd/mengasr-worker.service
[Unit]
Description=MengASR2 Worker — ASR 模型推理进程
After=network.target
Wants=network.target

[Service]
Type=simple
User=bchnm
WorkingDirectory=/srv/mengasr/server

# ── 配置文件 ──────────────────────────────────────────────
Environment=MENGASR_CONFIG=/srv/mengasr/mengasr.yaml
Environment=HF_ENDPOINT=https://hf-mirror.com
Environment=PIP_CONSTRAINT=/srv/mengasr/constraints.txt

ExecStart=/srv/mengasr/.venv/bin/python -m mengasr_server --worker
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mengasr-worker

# GPU 访问
SupplementaryGroups=video

# 资源限制
LimitNOFILE=65536
MemoryMax=24G

[Install]
WantedBy=multi-user.target
